# systemd/pi-agents-daemon.service
[Unit]
Description=pi Multi-Agent Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.pi/pi-commands/daemon
ExecStart=/usr/bin/node --experimental-strip-types %h/.pi/pi-commands/daemon/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
